# cloudtribe.service
[Unit]
Description=CloudTribe (Next.js)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/cloudtribe-2.0
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudtribe

# Environment variables if needed
Environment=NODE_ENV=production

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

